# roles/result/files/result.service
[Unit]
Description=result
After=network.target

[Service]
WorkingDirectory=/etc/appResult/result
ExecStart=npm start
Restart=always

[Install]
WantedBy=multi-user.target